# content/posts/2023-05-30-my-offline-music-streaming-setup/assets/gonic.service
[Install]
WantedBy=default.target

[Service]
ExecStart=gonic -config-path /etc/gonic/gonic.conf

[Unit]
After=network.target
After=sound.target
Description=Gonic server
Documentation=https://github.com/sentriz/gonic
